# supervisor/reset.ini
[program:admama_reset]
command=/usr/bin/php bin/console ADMAMA:endless_command RESET
directory=/data/http/AD-API/
user=root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/data/supervisor/log/reset.log
stderr_logfile=/data/supervisor/log/reset.log
