# supervisor/socket.ini
[program:admama_socket]
command=/usr/bin/php bin/console ADMAMA:link_to_socket
directory=/data/http/AD-API/
user=root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/data/supervisor/log/socket.log
stderr_logfile=/data/supervisor/log/socket.log
